# railsbase/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx
autorestart=true

[program:unicorn]
command=bundle exec unicorn -c config/unicorn.rb
autorestart=true
